# systemd/system/sync-package-database.service
[Unit]
Description=pacman -Sy

[Service]
ExecStart=/usr/sbin/pacman -Sy

[Install]
WantedBy=multi-user.target