# systemd/.config/systemd/user/sync-notes.timer
[Unit]
Description=Sync markdown notes daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
